# weather.cron
17 10 * * *  /home/ec2-user/DataScience/Stddev/run.sh
